# scripts/launchd/com.sdd.nightly-review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sdd.nightly-review</string>
    
    <key>ProgramArguments</key>
    <array>
        <!-- UPDATE THIS PATH to your project location -->
        <string>/Users/adrian/Documents/personal-repos/sdd-auto/scripts/nightly-review.sh</string>
    </array>
    
    <key>WorkingDirectory</key>
    <!-- UPDATE THIS PATH to your project location -->
    <string>/Users/adrian/Documents/personal-repos/sdd-auto</string>
    
    <!-- Run at 10:30 PM every day -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>22</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>
    
    <key>StandardOutPath</key>
    <!-- UPDATE THIS PATH to your project location -->
    <string>/Users/adrian/Documents/personal-repos/sdd-auto/logs/nightly-review.log</string>
    
    <key>StandardErrorPath</key>
    <!-- UPDATE THIS PATH to your project location -->
    <string>/Users/adrian/Documents/personal-repos/sdd-auto/logs/nightly-review.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <!-- UPDATE THIS to your home directory -->
        <string>/Users/adrian</string>
    </dict>
    
    <!-- Don't run if missed (e.g., laptop was asleep) -->
    <key>StartInterval</key>
    <false/>
</dict>
</plist>
